# ansible/roles/system/files/xhc.service
[Unit]
Description=Razer Blade XHC
After=syslog.target network.target

[Service]
ExecStart=echo XHC | tee /proc/acpi/wakeup
